# systemd/run-dm.service
[Unit]
Description=run dm

[Service]
ExecStart=/usr/bin/systemctl start ly

[Install]
WantedBy=multi-user.target
